# package/debian2/DEBIAN/prerm
#!/bin/bash
systemctl stop things_plex_service_template.service
systemctl disable things_plex_service_template.service
