# package/debian1/DEBIAN/postinst
#!/bin/bash

set -e

# summary of how this script can be called:
#		* <postinst> `configure' <most-recently-configured-version>
#		* <old-postinst> `abort-upgrade' <new version>
#		* <conflictor's-postinst> `abort-remove' `in-favour' <package>
#		  <new-version>
#		* <postinst> `abort-remove'
#		* <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#		  <failed-install-package> <version> `removing'
#		  <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

add_user_and_group()
{
	if ! getent group thingsplex >/dev/null
	then
		echo -n "Adding group thingsplex.."
		addgroup --quiet \
			--system \
			thingsplex 2>/dev/null || true
		echo "..done"
	fi

	if ! getent passwd conbee-ad >/dev/null
	then
		echo -n "Adding user conbee-ad.."
		adduser --quiet \
			--system \
			--ingroup thingsplex \
			--no-create-home \
			--disabled-password \
			conbee-ad 2>/dev/null || true
		echo "..done"
	fi

}

set_ownership()
{
	chown -R conbee-ad:thingsplex "/opt/thingsplex/conbee-ad"
	chmod 764 /var/log/thingsplex/conbee-ad || true
	chown conbee-ad:thingsplex /var/log/thingsplex/conbee-ad || true
}


case "$1" in
	configure)
	    echo "Configuring.."
		add_user_and_group
		set_ownership
		# Start conbee-ad
		if [ -f /opt/thingsplex/conbee-ad/conbee-ad ]
		then
		    systemctl enable conbee-ad.service
			systemctl start conbee-ad.service
		fi
	;;

	abort-upgrade|abort-remove|abort-deconfigure)
	;;

	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 1
	;;
esac



# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0