# package/debian1/DEBIAN/prerm
#!/bin/bash
systemctl stop conbee-ad.service
systemctl disable conbee-ad.service
