# package/debian/DEBIAN/prerm
#!/bin/bash
systemctl stop conbee.service || true
systemctl disable conbee.service || true
