# package/debian/DEBIAN/prerm
#!/bin/bash
systemctl stop conbee-ad.service || true
systemctl disable conbee-ad.service || true
